# .config/autostart/bitcoin.desktop
[Desktop Entry]
Type=Application
Name=Bitcoin
Exec=/usr/bin/bitcoin-qt -min
Terminal=false
Hidden=false
